# roles/spads/templates/clear_rapid_packages.timer.j2
[Unit]
Description=Clear old rapid packages daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
