# roles/spads/templates/clear_logs_and_replays.timer.j2
[Unit]
Description=Clear SPADS logs and replays everyday at midnight

[Timer]
OnActiveSec=5sec
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
